# bitool.service
[Unit]
Description=BiTool - Professional Email & Facebook Tools Application
After=network.target

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/home/root/webapp
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/local/bin/gunicorn --config gunicorn_config.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
